# fpm/fpm-8.3/supervisor/conf.d/laravel-reverb.conf
# /etc/supervisor/conf.d/laravel-reverb.conf
# Конфигурация для Laravel Reverb Server

[program:laravel-reverb]
process_name=%(program_name)s
command=php /data/esia.local/backend/artisan reverb:start --host=0.0.0.0 --port=8080
directory=/data/esia.local/backend
autostart=true
autorestart=true
startsecs=1
startretries=3
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/laravel-reverb.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=2
stopwaitsecs=30
killasgroup=true
priority=998