# fpm/fpm-8.3/supervisor/conf.d/laravel-schedule.conf

# /etc/supervisor/conf.d/laravel-schedule.conf
# Конфигурация для Laravel Scheduler (опционально)

[program:laravel-schedule]
process_name=%(program_name)s
command=php /data/esia.local/backend/artisan schedule:work
directory=/data/esia.local/backend
autostart=true
autorestart=true
startsecs=1
startretries=3
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/laravel-schedule.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=2
stopwaitsecs=30
killasgroup=true
priority=997